# themes/luci-theme-bosminer-plus/root/etc/uci-defaults/30_luci-theme-miner
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap=/luci-static/miner
	set luci.main.mediaurlbase=/luci-static/miner
	commit luci
EOF
exit 0
